# tutorials/setup-rdbox-hq/build/roles/install_softethervpn/files/vpnbridge.service
[Unit]
Description=SoftEtherVPN(bridge)
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/vpnbridge/vpnbridge start
ExecStop=/usr/local/vpnbridge/vpnbridge stop

[Install]
WantedBy=multi-user.target

